# src/life_world_model/daemon/com.lwm.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lwm.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/lwm</string>
        <string>watch</string>
        <string>--interval</string>
        <string>60</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/lwm-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/lwm-daemon.err</string>
</dict>
</plist>
